# crontab
@daily /home/nas/archilab-backup/backup-docker-db prox-prod prox-project-service project-db project-service project-db
